# board/avnet/zedboard/overlay/etc/profile.d/zedboard_prompt.sh
export PS1="\[\e[0:32m\]zed\[\e[m\]\[\e[1;32m\]board\[\e[m\]\[\e[0;32m\]>\[\e[m\] "
